# conf/crontab
# Run phonehome script every minute
* * * * * /home/pi/bin/phonehome
