# pkgs/yodelgw/DEBIAN/postrm
#!/bin/sh
set -e

systemctl disable lora-gateway-bridge.service
systemctl disable packet-forwarder.service
delgroup --system sdr
deluser --system lora

exit 0
