# pkgs/yodelgw/DEBIAN/postrm
#!/bin/sh
set -e

systemctl daemon-reload
systemctl disable dhcpcd.service
systemctl enable dhcpcd.service

delgroup --system sdr
deluser --system lora

exit 0
