# pkgs/yodelgw/DEBIAN/postinst
#!/bin/sh

# These need to succeed
addgroup --system sdr
adduser --system --group --no-create-home --home /nonexistent lora
adduser lora gpio
adduser lora spi
adduser lora dialout

systemctl daemon-reload
systemctl disable dhcpcd.service
systemctl enable dhcpcd.service

systemctl enable lora-gateway-bridge.service
systemctl enable packet-forwarder.service
systemctl start lora-gateway-bridge.service
systemctl start packet-forwarder.service

udevadm control --reload-rules

exit 0
