# pkgs/yodelgw/lib/systemd/system/dhcpcd.service.d/waitfix.conf
# This actually makes systemd's target network-online wait for dhcpcd to get an IP.
#
# Craig Hesling <craig@hesling.com>
# August 13, 2016

[Unit]
Before=network-online.target

[Install]
WantedBy=network-online.target
